# deployment/uwsgi/scitex_cloud_prod.service
[Unit]
Description=uWSGI instance to serve SciTeX Cloud
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/scitex-cloud
Environment="PATH=/var/www/scitex-cloud/venv/bin"
EnvironmentFile=/var/www/scitex-cloud/.env
ExecStart=/var/www/scitex-cloud/venv/bin/uwsgi --ini /etc/uwsgi/scitex_cloud.ini

# Restart policy
Restart=always
RestartSec=10s

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
