# deployment/.archive/manual/systemd/scitex-cleanup-jobs.timer
[Unit]
Description=SciTeX Scholar - Periodic cleanup of stale jobs (every 5 minutes)
Requires=scitex-cleanup-jobs.service

[Timer]
# Run every 5 minutes to catch stale jobs quickly
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
